# systemd/admin-backup.service
[Unit]
Description=Admin backup job
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/srv/admin/env/backup.env
ExecStart=/opt/homelab-admin-node/scripts/backup.sh
