# systemd/admin-converge.timer
[Unit]
Description=Periodic admin convergence

[Timer]
OnBootSec=5m
OnUnitActiveSec=30m
Unit=admin-converge.service

[Install]
WantedBy=timers.target
