# systemd/admin-backup.timer
[Unit]
Description=Daily admin backup

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
Unit=admin-backup.service

[Install]
WantedBy=timers.target
