# systemd/admin-unlock.path
[Unit]
Description=Watch for age key installation

[Path]
PathChanged=/etc/sops/age/keys.txt
Unit=admin-converge.service

[Install]
WantedBy=multi-user.target
