# systemd/admin-stack@.service
[Unit]
Description=Admin stack %i
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=true
WorkingDirectory=/srv/admin/stacks/%i
EnvironmentFile=-/srv/admin/env/%i.env
ExecStart=/usr/bin/docker compose -f compose.yaml up --wait --wait-timeout 300
ExecStop=/usr/bin/docker compose -f compose.yaml down

[Install]
WantedBy=admin-stacks.target
